# hc-stack-install/files/vault-agent-watcher.path
[Path]
PathChanged=/etc/consul.d/vault_token

[Install]
WantedBy=multi-user.target