# ansible/roles/hc_stack_install/files/consul-watcher.path
[Path]
PathModified=/etc/consul.d/consul.hcl

[Install]
WantedBy=multi-user.target