# ansible/hc-stack-install/files/nomad-watcher.path
[Path]
PathModified=/etc/nomad.d/nomad.hcl

[Install]
WantedBy=multi-user.target